# system_files/shared/usr/lib/systemd/system/caracal-fix-sddm-background.service
[Unit]
Description=Repair stale SDDM greeter background override
DefaultDependencies=no
Wants=usr-share-sddm-themes.mount
After=local-fs.target systemd-tmpfiles-setup.service usr-share-sddm-themes.mount
Before=sddm.service
ConditionPathExists=/usr/share/sddm/themes/caracal/theme.conf

[Service]
Type=oneshot
ExecStart=/usr/libexec/caracal-fix-sddm-background

[Install]
WantedBy=multi-user.target
