# system_files/shared/usr/lib/systemd/system/usr-share-sddm-themes.mount
# Workaround to allow KDE Discover to install SDDM themes on an immutable system.
# Without this overlay, /usr/share/sddm/themes is read-only and Discover cannot
# write there. The overlay redirects writes to /var/sddm_themes/ which persists.
# TODO: Remove when sddm supports installing themes outside /usr/share.
# See https://github.com/sddm/sddm/issues/1561

[Unit]
Description=KDE writable SDDM themes workaround
RequiresMountsFor=/usr /var
ConditionPathExists=/usr/share/sddm
ConditionPathExists=/var/sddm_themes/themes
ConditionPathExists=/var/sddm_themes/themes.work

[Mount]
Type=overlay
What=overlay
Where=/usr/share/sddm/themes
Options=lowerdir=/usr/share/sddm/themes,upperdir=/var/sddm_themes/themes,workdir=/var/sddm_themes/themes.work

[Install]
WantedBy=multi-user.target
